# hardbreaks.service
# place in  /etc/systemd/system/hardbreaks.service
[Unit]
Description=HardBreaks LED Matrix - run animation
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/HardBreaks
ExecStart=/usr/bin/node /home/pi/HardBreaks/pi/run.js

# wait 10 seconds if crashes, then restart
Restart=always
RestartSec=10

# Restart every 1 hour 1 minute and 3 seconds
RuntimeMaxSec=3663

[Install]
WantedBy=multi-user.target

